# backend/sni-router-ui.service
[Unit]
Description=SNI-Router UI backend
After=network.target

[Service]
WorkingDirectory=/opt/sni-router-ui/backend
# Host / port / DB path come from the local config file (editable in the UI).
EnvironmentFile=/etc/sni-router-ui/ui.conf
ExecStart=/opt/sni-router-ui/backend/.venv/bin/uvicorn app:app --host ${SNI_UI_HOST} --port ${SNI_UI_PORT}
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
