# backend/sni-router-ui.service
[Unit]
Description=SNI-Router UI backend
After=network.target

[Service]
WorkingDirectory=/opt/sni-router-ui/backend
Environment=SNI_UI_DB=/opt/sni-router-ui/data/sni-ui.db
Environment=SNI_UI_HOST=127.0.0.1
Environment=SNI_UI_PORT=8899
ExecStart=/opt/sni-router-ui/backend/.venv/bin/uvicorn app:app --host ${SNI_UI_HOST} --port ${SNI_UI_PORT}
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
