# agent/sni-router-agent.service
[Unit]
Description=SNI-Router UI metrics agent
After=network.target

[Service]
# Token / bind / port come from the local config file (written by the installer).
EnvironmentFile=/etc/sni-router-agent/agent.conf
ExecStart=/usr/bin/python3 /opt/sni-router-agent/agent.py
Restart=on-failure
RestartSec=3
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes

[Install]
WantedBy=multi-user.target
